# config/systemd/user/gui.target
[Unit]
Description=GUI
Requires=default.target
After=default.target
