# config/systemd/user/polybar@.service
[Unit]
Description=Lightweight Status Bar
Documentation=man:polybar(1)

[Service]
ExecStart=/usr/bin/polybar -r %i

[Install]
WantedBy=x11.target
