# config/systemd/user/x11.target
[Unit]
Description=X11 GUI
Requires=gui.target
After=gui.target
