# tag-systemd/config/systemd/user/compton.service
[Unit]
Description=Simple X11 Compositor
Documentation=man:compton(1)

[Service]
ExecStart=/usr/bin/compton

[Install]
WantedBy=x11.target
