# config/systemd/dunst.service
[Unit]
Description=Lightweight Notification Daemon
Documentation=man:dunst(1)

[Service]
ExecStart=/usr/bin/dunst

[Install]
WantedBy=x11.target
