# config/systemd/sound.target
[Unit]
Description  = Sound Server
Requires     = dbus.socket
After        = dbus.socket

[Install]
Alias = default.target
