# config/systemd/xbindkeys.service
[Unit]
Description=X11 Key Grabber
Documentation=man:xbindkeys(1)

[Service]
ExecStart=/usr/bin/xbindkeys -n

[Install]
WantedBy=x11.target
