# config/systemd/xorg.target
[Unit]
Description = X Window System

[Install]
Alias = default.target
