# home/private_dot_config/systemd/user/windows-ssh-proxy.socket.tmpl
{{- if .wsl -}}
[Unit]
Description=Windows SSH Agent Proxy Socket

[Socket]
ListenStream=%t/windows-ssh-proxy.sock
Accept=yes

[Install]
WantedBy=sockets.target
{{- end -}}
